# data/systemd/netmntd.service
[Unit]
Description=netmnt network mount daemon
Documentation=https://github.com/
# D-Bus activated; started on first call to org.netmnt.
After=network-online.target
Wants=network-online.target

[Service]
Type=dbus
BusName=org.netmnt
ExecStart=/usr/bin/netmntd
# Hardening (mount needs CAP_SYS_ADMIN, so we cannot drop everything).
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_OVERRIDE CAP_SETUID CAP_SETGID
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/etc/systemd/system /run /home
NoNewPrivileges=false
Restart=on-failure

[Install]
WantedBy=multi-user.target
