# data/systemd/netmntd.service
[Unit]
Description=netmnt network mount daemon
Documentation=https://github.com/thongor77/netmnt
After=network-online.target
Wants=network-online.target

[Service]
Type=dbus
BusName=org.netmnt
ExecStart=/usr/bin/netmntd
Restart=on-failure

# IMPORTANT: the daemon performs mounts on behalf of the host. With a private
# mount namespace (the default once any Protect*/sandboxing option is set) those
# mounts would NOT propagate to user sessions. MountFlags=shared keeps mount
# propagation working. Any future hardening must preserve this.
MountFlags=shared

[Install]
WantedBy=multi-user.target
